# fs/etc/supervisor.d/mysql.ini
[program:mysql]
command = /usr/sbin/mysqld -u mysql
pidfile = /run/mysqld/mysqld.pid
user = mysql
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0