# fs/etc/supervisor.d/varnish.ini
[program:varnish]
command = /usr/sbin/varnishd -f /etc/varnish/default.vcl -a 0.0.0.0:81 -s malloc,1G -F
user = root
autostart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0